# systemd/ha-tux-smart.timer
[Unit]
Description=Daily SMART health scan for ha-tux

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
