# systemd/ha-tux-smart.service
[Unit]
Description=Scan SMART health for ha-tux

[Service]
Type=oneshot
TimeoutStartSec=120
ExecStart=/opt/ha-tux/venv/bin/ha-tux-smart-scan
